# btrfs-snap/btrfs-daily-snapshot.service
[Unit]
Description=Snapshot filesystem daily

[Service]
User=root
Type=simple
ExecStart=/root/bin/btrfs-snap date
