# btrfs-snap/btrfs-daily-snapshot.timer
[Unit]
Description=Snapshot filesystem daily

[Timer]
OnUnitActiveSec=24h
Unit=btrfs-daily-snapshot.service

[Install]
WantedBy=multi-user.target
